# scripts/claude-traffic-guard.timer
[Unit]
Description=Claude Traffic Guard Timer - Check exit IP every 30s

[Timer]
OnBootSec=30
OnUnitActiveSec=30
AccuracySec=5

[Install]
WantedBy=timers.target
