# scripts/claude-traffic-guard.service
[Unit]
Description=Claude Traffic Guard - Block Claude API if exit IP is wrong
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
# 部署时改成你的实际路径，如 /home/<user>/happyclaw/scripts/claude-traffic-guard.sh
ExecStart=/path/to/happyclaw/scripts/claude-traffic-guard.sh
User=root
